# waterscribe.service
[Unit]
Description=WaterScribe Web Application
After=network.target

[Service]
Type=simple
User=rcampbell
WorkingDirectory=/home/rcampbellyy/waterscribe
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
